# scripts/systemd/openclaw-audit-report.timer
[Unit]
Description=Run OpenClaw security audit every 10 minutes

[Timer]
OnBootSec=5m
OnUnitActiveSec=10m
Persistent=true
Unit=openclaw-audit-report.service

[Install]
WantedBy=timers.target
