# scripts/systemd/openclaw-audit-report.service
[Unit]
Description=Run OpenClaw security audit and upload result
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=%h/.config/openclaw-audit.env
ExecStart=%h/.local/bin/openclaw-audit-reporter.sh
